# systemd/xvfb@.service
[Unit]
Description=Virtual Framebuffer X Server (:%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
EnvironmentFile=/etc/sandboxmcp/visual.env
Environment=DISPLAY=:%i
ExecStart=/usr/bin/Xvfb :%i -screen 0 1920x1080x24
Restart=always

[Install]
WantedBy=multi-user.target
